# defaults/private_Library/LaunchAgents/com.sebastienrousseau.corralctl.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.sebastienrousseau.corralctl</string>

	<key>ProgramArguments</key>
	<array>
		<string>{{ .chezmoi.homeDir }}/.local/bin/corralctl-sync.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>0</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>RunAtLoad</key>
	<false/>

	<key>ProcessType</key>
	<string>Background</string>

	<key>StandardOutPath</key>
	<string>{{ .chezmoi.homeDir }}/Library/Logs/corralctl.launchd.log</string>
	<key>StandardErrorPath</key>
	<string>{{ .chezmoi.homeDir }}/Library/Logs/corralctl.launchd.log</string>
</dict>
</plist>
